# crontab
* * * * * php /var/www/artisan schedule:run >> 2>&1
# Leave one blank line below or cron will fail
